# network/deploy/v86net-gateway.service
[Unit]
Description=VM network gateway
After=network-online.target v86net-network.service
Wants=network-online.target
Requires=v86net-network.service

[Service]
Type=simple
ExecStart=/usr/local/bin/v86net-gateway -listen 127.0.0.1:8086 -tap v86tap0 -admin-token-file /etc/v86net/admin.token -allow-origin https://avirads.github.io
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# Opening an existing TAP requires network-administration capability.
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
